# cron/crontab
0 * * * * root bash /app/cron/scripts/delete_videos.sh
